# systemd/8rupia-web.service
[Unit]
Description=8Rupia Next.js web server (port 3002 — internal only)
After=network.target docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/8rupia
EnvironmentFile=/var/www/8rupia/.env.production
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3002
ExecStart=/usr/bin/node node_modules/.bin/next start -H 127.0.0.1 -p 3002
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
